# etc/rc.local
#!/bin/sh -e

stty -F /dev/ttyAMA0 115200 cs8 -cstopb -parenb -brkint -imaxbel -opost -onlcr -isig -icanon -iexten -echo -echoe -echok -echoctl -echoke
modprobe xgo-drv
st7789_server
chgrp video /tmp/lcd0
chmod g+w /tmp/lcd0

exit 0
